# etc/systemd/system/socat.service
[Unit]
Description=Forward TCP from 80 to 8888 (bidirectional)
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/socat TCP-LISTEN:80,reuseaddr,fork TCP:8888
Restart=always

[Install]
WantedBy=multi-user.target
